# systemd_example/amdctl-boot.service
[Unit]
Description=Run amdctl after boot

[Service]
Type=oneshot
ExecStart=/usr/sbin/amdctl.sh
